# scripts/onesibox.service
[Unit]
Description=OnesiBox Client
After=network-online.target graphical.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=debian
ExecStart=/usr/bin/node /opt/onesibox/src/main.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
# Display environment - supports both X11 and Wayland
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/opt/onesibox

# Sicurezza
# NOTA: NoNewPrivileges=true non puo' essere usato perche' il servizio
# ha bisogno di sudo per comandi di reboot/restart (via /etc/sudoers.d/onesibox)
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/onesibox /var/log/onesibox /home/debian
PrivateTmp=true

# Sensitive configuration via environment variables (recommended for production)
# Uncomment and set these instead of storing in config.json:
# Environment=ONESIBOX_SERVER_URL=https://onesiforo.example.com
# Environment=ONESIBOX_APPLIANCE_ID=your-appliance-uuid
# Environment=ONESIBOX_TOKEN=your-secret-token
#
# Or load from a separate credentials file (more secure):
# EnvironmentFile=/opt/onesibox/config/credentials.env

[Install]
WantedBy=graphical.target
