# scripts/onesibox.service
[Unit]
Description=OnesiBox Client
After=network-online.target graphical.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
# User will be set dynamically by install.sh based on detected kiosk user
User=admin
ExecStart=/usr/bin/node /opt/onesibox/src/main.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
# Display environment - supports both X11 and Wayland
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
# XDG_RUNTIME_DIR will be set dynamically based on user UID
Environment=XDG_RUNTIME_DIR=/run/user/1000
WorkingDirectory=/opt/onesibox

# Sicurezza
# ProtectHome=read-only to allow access to Playwright browser cache
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/opt/onesibox /var/log/onesibox /home/admin/.onesibox-zoom /home/admin/.config/chromium /home/admin/.cache/ms-playwright
PrivateTmp=true

[Install]
WantedBy=graphical.target
